# docker/root/etc/services.d/kokorobot/run
#!/bin/sh

cd /kokorobot
bin/hubot -a "$KOKOROBOT_HUBOT_ADAPTER" -n "$KOKOROBOT_HUBOT_NAME" mohemohe-kokorobot
